# projects/Rockchip/devices/RK322X/packages/irq-balancer/system.d/irq-balancer.service
[Unit]
Description=Configure irq
Before=kodi.service
After=network-online.target graphical.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/irq-balancer

[Install]
WantedBy=kodi.target
